# etc/supervisor/conf.d/.conf
[program:bms]
command=/path/to/myprogram --option arg
autostart=true
autorestart=true
stderr_logfile=/var/log/myprogram.err.log
stdout_logfile=/var/log/myprogram.out.log